# backup.timer
[Unit]
Description=Run Blacket database backup every 12 hours

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target